# cli/php-webserver.conf
[program:php-cli-example]
command=php /home/lorna/scripts/cli/supervise-me.php
autostart=false
autorestart=true

